# scripts/jayu-systemd.timer
[Unit]
Description=Run Jayu stock automation every 4 hours

[Timer]
OnCalendar=*-*-* 00,04,08,12,16,20:00:00
Persistent=true
Unit=jayu-systemd.service

[Install]
WantedBy=timers.target
